# installdir/var/spool/cron/crontabs/root
0 	1 	* 	* 	*	ntpdate 133.100.11.8 && hwclock -w
